# deploy/systemd/peacock-socket.service
[Unit]
Description=Peacock Time Socket Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/peacock-time-v2
Environment=NODE_ENV=production
Environment=SOCKET_PORT=3001
Environment=SOCKET_SECRET=replace-with-secret
ExecStart=/usr/bin/node server/socket-server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
